# packaging/service/dns2doh.service
[Unit]
Description=DNS to DoH
Wants=network-online.target
After=network-online.target
[Service]
ExecStart=/usr/share/dns2doh/dns2doh
ExecStop=
Restart=on-failure
RestartSec=5
TimeoutStartSec=60
TimeoutStopSec=5
KillMode=process
Type=forking
[Install]
WantedBy=multi-user.target